# deployment/make_s1_ard/cron/make_s1_ard.cron
# Run Every Day at 01:00
0 1 * * * /usr/bin/make_s1_ard