# deployment/get_s1_raw/cron/get_s1_raw.cron
# Run Every Day at 00:00
0 0 * * * /usr/local/bin/get_s1_raw